# packaging/systemd/openmailstack.service
[Unit]
Description=OpenMailStack Webmail Backend
After=network.target mariadb.service dovecot.service postfix.service

[Service]
Type=simple
User=openmailstack
Group=openmailstack
WorkingDirectory=/opt/openmailstack-backend
Environment=NODE_ENV=production
EnvironmentFile=/etc/openmailstack/webmail-backend.env
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5
PrivateTmp=true
ProtectSystem=full
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
